# deploy/native/linux/tunnara-server.service
[Unit]
Description=Tunnara Platform Server 0.2
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tunnara
Group=tunnara
EnvironmentFile=-/etc/tunnara/server.env
ExecStart=/usr/bin/node /opt/tunnara/runtime/node/bin/tunnara-server.mjs serve-all --data-dir /var/lib/tunnara
WorkingDirectory=/opt/tunnara
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tunnara
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
