# deploy/standalone/linux/tunnara-agent.service
[Unit]
Description=Tunnara Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tunnara-agent
Group=tunnara-agent
Environment=TUNNARA_CONFIG_DIR=/var/lib/tunnara-agent
ExecStart=/usr/local/bin/tunnara serve --config-dir /var/lib/tunnara-agent
WorkingDirectory=/var/lib/tunnara-agent
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tunnara-agent
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
